# revolution@.service.in
[Unit]
Description=Blue Sky Solar Racing Revolution %I

[Service]
Type=simple
User=@username@
ExecStart=@bindir@/%I
ExecStop=@bindir@/client client %I EXIT
Restart=on-failure

[Install]
WantedBy=multi-user.target
